# install/direwolf.service
[Unit]
Description=Direwolf KISS TNC for Meshtastic Bridge
Documentation=https://github.com/wb2osz/direwolf
After=sound.target
Before=msh-bridge.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/msh-bridge
ExecStart=/usr/bin/direwolf -c /opt/msh-bridge/direwolf.conf -t 0

# Restart policy
Restart=always
RestartSec=5

# Resource limits
MemoryMax=256M
CPUQuota=50%

# Logging to /tmp (tmpfs in RAM) to reduce SD card wear
StandardOutput=append:/tmp/direwolf.log
StandardError=append:/tmp/direwolf.log
SyslogIdentifier=direwolf

[Install]
WantedBy=multi-user.target
